# app/dish-monitor.timer
[Unit]
Description=Run dish monitor every 10 seconds

[Timer]
OnBootSec=10s
OnUnitActiveSec=10s
AccuracySec=1s

[Install]
WantedBy=timers.target
